# liquidcfg_init.service
[Unit]
Description=Initialize liquidctl

[Service]
Type=oneshot
ExecStart=liquidctl initialize all

[Install]
WantedBy=default.target
